# systemd/hauptbahnhof.target
[Unit]
Description=Hauptbahnhof Hackerspace Control System
Wants=mosquitto.service

[Install]
WantedBy=multi-user.target
