# systemd/hauptbahnhof.target
[Unit]
Description=Hauptbahnhof Hackerspace Control System

[Install]
WantedBy=multi-user.target